# pack/systemd/aevum-sysctl-apply.service
[Unit]
Description=Apply Aevum sysctl hardening
ConditionPathExists=/var/lib/aevum/workstation/identity/device_ed25519_sk.pem
After=local-fs.target
Wants=local-fs.target

[Service]
Type=oneshot
ExecStart=/sbin/sysctl --system





# --- Aevum service hardening (control-plane, non-gating) ---
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
