# pack/systemd/aevum-timechain.service
# /etc/systemd/system/aevum-timechain.service
[Unit]
Description=Aevum TimeChain Daemon (1 block/sec)
ConditionPathExists=/var/lib/aevum/workstation/identity/device_ed25519_sk.pem
After=aevum-identity-bootstrap.service aevum-tpm-unseal.service
Wants=aevum-identity-bootstrap.service aevum-tpm-unseal.service

[Service]
Type=simple
RuntimeDirectory=aevum
RuntimeDirectoryMode=0755
ExecStart=/usr/local/sbin/aevum_timechain_daemon.py --base /var/lib/aevum
Restart=always
RestartSec=1
User=root
Group=root
UMask=0077





# --- Aevum service hardening (data-plane, non-gating) ---
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
SystemCallArchitectures=native
UMask=0077
ReadWritePaths=/var/lib/aevum /run/aevum /var/log /etc/aevum

[Install]
WantedBy=multi-user.target
