# pack/systemd/aevum-segment.service
[Unit]
Description=Aevum segment/rotate receipt logs (best-effort)
ConditionPathExists=/var/lib/aevum/workstation/identity/device_ed25519_sk.pem
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/aevum_segment.py --base /var/lib/aevum/workstation
Nice=10
IOSchedulingClass=best-effort
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
