# pack/systemd/aevum-apt-capture.service
[Unit]
Description=Aevum APT Supply-Chain Capture
ConditionPathExists=/var/lib/aevum/workstation/identity/device_ed25519_sk.pem
After=local-fs.target
Wants=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/aevum_apt_capture.py --base /var/lib/aevum --instance workstation

NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
SystemCallArchitectures=native
UMask=0077
