# pack/systemd/aevum-segment.timer
[Unit]
Description=Run Aevum receipt segmentation periodically

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
