# pack/systemd/aevum-apt-capture.timer
[Unit]
Description=Aevum APT Supply-Chain Capture (daily)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
