# systemd/aevum-quarantine.target
[Unit]
Description=Aevum Quarantine Target (unbound appliance mode)
Requires=basic.target
After=basic.target
AllowIsolate=yes
